# supervisor.conf
[program:python-podcast]
command=docker-compose -f production.yml up
directory=/home/python-podcast/python-podcast
redirect_stderr=true
autostart=true
autorestart=true
priority=10
